# kafka/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh --override zookeeper.connect=zookeeper:2181
autostart=true
autorestart=true